# deploy/systemd/megatrader-deploy.service
[Unit]
Description=MegaTrader Deploy Service

[Service]
Type=oneshot
ExecStart=/opt/megatrader/deploy.sh
User=root
