# deploy/systemd/megatrader-webhook.service
[Unit]
Description=MegaTrader Webhook Service
After=network.target

[Service]
User=deploy
WorkingDirectory=/var/www/MegaTrader/deploy
ExecStart=/opt/megatrader/venv/bin/uvicorn webhook_server:app --host 0.0.0.0 --port 9000
Restart=always

[Install]
WantedBy=multi-user.target
